# ajedrez/supervisor.conf
[program:chess]
command=/bin/java -jar /root/serverajedrez.jar
autostart=true
autorestart=true